# config/autostart/nvidia-force-full-composition.desktop
[Desktop Entry]
Type=Application
Name=nvidia-force-full-compositon
Exec=/home/tarcisio/scripts/nvidia-force-full-composition.sh
